# supervisord/queues/simple-queue.conf
[program:ampq-consume-simple-queue]
command=drush amqp:consume simple-queue
numprocs=1
autostart=true
autorestart=true
